# apkovl/etc/init.d/takeover
#!/sbin/openrc-run

depends() {
	after displayservice
}

pidfile="/run/${RC_SVCNAME}.pid"

start() {
    /opt/tinybox/takeover.sh &> /dev/console &
}
